# smf/smf.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!-- 
    Manifest automatically generated by smfgen.
 -->
<service_bundle type="manifest" name="application-image-networking" >
    <service name="application/image-networking" type="service" version="1" >
        <create_default_instance enabled="true" />
        <dependency name="dep0" grouping="require_all" restart_on="error" type="service" >
            <service_fmri value="svc:/milestone/multi-user:default" />
        </dependency>
        <exec_method type="method" name="start" exec="/opt/custom/smf/imgnet.sh start &amp;" timeout_seconds="10" >
            <method_context >
                <method_credential user="root" group="root" />
            </method_context>
        </exec_method>
        <exec_method type="method" name="stop" exec="/opt/custom/smf/imgnet.sh stop" timeout_seconds="30" >
            <method_context >
                <method_credential user="root" group="root" />
            </method_context>
        </exec_method>
        <template >
            <common_name >
                <loctext xml:lang="C" >Configure networking for bhyve guest image creation</loctext>
            </common_name>
        </template>
    </service>
</service_bundle>
